# ops/apex_v4.service
[Unit]
Description=APEX V4 Algorithmic Trading System
After=network.target redis.service postgresql.service
Requires=redis.service postgresql.service

[Service]
Type=simple
User=apex
WorkingDirectory=/opt/apex_v4
ExecStart=/opt/apex_v4/venv/bin/python -m src.pipeline
Restart=always
RestartSec=5
Environment=PYTHONPATH=/opt/apex_v4
EnvironmentFile=/opt/apex_v4/config/secrets.env

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=apex_v4

# Resource limits
LimitNOFILE=65535

[Install]
WantedBy=multi-user.target
